# django__api/crontab
*/1 * * * * /process_task_script.sh >> /var/log/cron.log
